# backup/backup.timer
[Unit]
Description=rsync backup

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
